# gpu-switcher.service
[Unit]
Description=GPU switching service

[Service]
Type=oneshot
ExecStart=/usr/bin/gpu-switcher 

[Install]
WantedBy=multi-user.target
